# t.conf
[program:t]
command=python3 /home/curator2/t/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/t/test.err.log
stdout_logfile=/home/curator2/t/test.out.log